# fy26_pmo/com.openclaw.fy26-pmo-report.template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.fy26-pmo-report</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/admin/.openclaw/workspace/fy26_pmo/run.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/admin/.openclaw/workspace/fy26_pmo</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/admin</string>
        <!-- ⚠️ SECURITY: Set these in your local environment, not here! -->
        <!-- JIRA_API_TOKEN will be loaded from ~/.openclaw/workspace/.jira-config -->
        <!-- QQ_MAIL_PASSWORD will be loaded from ~/.openclaw/workspace/.jira-config -->
    </dict>
    
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/admin/.openclaw/workspace/fy26_pmo/logs/launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/admin/.openclaw/workspace/fy26_pmo/logs/launchd_error.log</string>
</dict>
</plist>
