# projects/cntin730-report/config/com.openclaw.cntin730-report.template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.cntin730-report</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/admin/.openclaw/workspace/projects/cntin730-report/scripts/run.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/admin/.openclaw/workspace/projects/cntin730-report</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/admin</string>
        <!-- ⚠️ SECURITY: Add these locally, DO NOT commit to git! -->
        <!--
        <key>JIRA_API_TOKEN</key>
        <string>YOUR_JIRA_API_TOKEN_HERE</string>
        <key>QQ_MAIL_PASSWORD</key>
        <string>YOUR_QQ_AUTH_CODE_HERE</string>
        -->
    </dict>
    
    <!-- 每个工作日12:00运行 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/admin/.openclaw/workspace/projects/cntin730-report/logs/cron.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/admin/.openclaw/workspace/projects/cntin730-report/logs/cron_error.log</string>
</dict>
</plist>
